# .config/systemd/user/applications.target
[Unit]
Description=Standard applications I always start
Requires=graphical-session.target
After=graphical-session.target
AllowIsolate=no