# .config/systemd/user/blueman-applet.service
[Unit]
Description=blueman applet

[Install]
WantedBy=graphical-session.target

[Service]
ExecStart=/usr/bin/blueman-applet