# .config/systemd/user/chromium.service
[Unit]
Description=Chromium Web Browser

[Service]
#Type=notify
ExecStart=/usr/bin/chromium

[Install]
WantedBy=applications.target
