# .config/systemd/user/fbautostart.service
[Unit]
Description=fbautostart

[Install]
WantedBy=xmonad.service

[Service]
Type=oneshot
ExecStart=/usr/bin/fbautostart