# .config/systemd/user/flameshot.service
[Unit]
Description=Flameshot Screenshot Tool

[Install]
WantedBy=graphical-session.target

[Service]
ExecStart=/usr/bin/flameshot
