# .config/systemd/user/nm-applet.service
[Unit]
Description=NetworkManager applet

[Install]
WantedBy=graphical-session.target

[Service]
ExecStart=/usr/bin/nm-applet
