# .config/systemd/user/parcellite.service
[Unit]
Description=ClipIt GTK Clipboard manager

[Install]
WantedBy=graphical-session.target

[Service]
ExecStart=/usr/bin/parcellite
