# .config/systemd/user/pasystray.service
[Unit]
Description=PulseAudio systray

[Install]
WantedBy=graphical-session.target

[Service]
ExecStart=/usr/bin/pasystray
