# .config/systemd/user/pavucontrol.service
[Unit]
Description=PulseAudio Volumecontrol

[Install]
WantedBy=applications.target

[Service]
ExecStart=%h/.local/state/nix/profile/bin/pavucontrol