# .config/systemd/user/status-notifier-watcher.service
[Unit]
Description=status-notifier-watcher

[Service]
ExecStart=/usr/bin/status-notifier-watcher