# .config/systemd/user/taffybar.service
[Unit]
Description=taffybar
Requires=xmonad.service
After=xmonad.service

[Install]
WantedBy=xmonad.service

[Service]
ExecStart=/usr/bin/taffybar