# .config/systemd/user/xmonad-systemd-session.target
[Unit]
Description=xmonad systemd session target
BindsTo=graphical-session.target
