# .config/systemd/user/xmonad.service
[Unit]
Description=xmonad window manager
After=graphical-session.target
Wants=graphical-session.target

[Service]
ExecStart=/usr/bin/xmonad