# .config/systemd/user/emacs.service.d/customdependency.conf
# https://wiki.archlinux.org/index.php/systemd#Examples
[Unit]
Requires=require-debian-package@elpa-systemd.service
After=require-debian-package@elpa-systemd.service