# .config/systemd/user/emacs.service.d/override.conf
[Service]
ExecStart=
ExecStart=%h/.nix-profile/bin/emacs --fg-daemon
ExecStop=
ExecStop=%h/.nix-profile/bin/emacsclient --eval "(kill-emacs)"