# .config/systemd/user/redshift.service.d/override.conf
[Install]
WantedBy=
WantedBy=graphical-session.target
